# seo-pipeline/com.seopipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.seopipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/macmini/dev/SEO/seo-pipeline/seo_pipeline.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>REPLACE_WITH_YOUR_KEY</string>
        <key>SERPER_API_KEY</key>
        <string>REPLACE_WITH_YOUR_SERPER_KEY</string>
        <key>GOOGLE_API_KEY</key>
        <string>REPLACE_WITH_YOUR_GOOGLE_API_KEY</string>
    </dict>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>/Users/macmini/dev/SEO/seo-pipeline/pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macmini/dev/SEO/seo-pipeline/pipeline.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/macmini/dev/SEO/seo-pipeline</string>
</dict>
</plist>
